# tmux/.config/systemd/user/tmux-session-save.service
[Unit]
Description=Save tmux session state

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tmux-session save
