# tmux/.config/systemd/user/tmux-session-save.timer
[Unit]
Description=Save tmux session state every 60 seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
